# scripts/scholartrace-mcp.service
[Unit]
Description=ScholarTrace MCP Server
After=network.target
# Secondary managed-deployment example; the main quick path is the tmux scripts in README.

[Service]
Type=simple
User=scholartrace
Group=scholartrace
WorkingDirectory=/opt/ScholarTrace
# Put SCHOLARTRACE_BIGMODEL_API_KEY and optional DeepXiv settings here.
EnvironmentFile=/etc/scholartrace/scholartrace.env
Environment=SCHOLARTRACE_MCP_TRANSPORT=sse
Environment=SCHOLARTRACE_MCP_HOST=0.0.0.0
Environment=SCHOLARTRACE_MCP_PORT=8001
Environment=SCHOLARTRACE_REMOTE_ACCESS_ENABLED=true
Environment=SCHOLARTRACE_ACCESS_TOKEN=g203-mcp
Environment=SCHOLARTRACE_DATA_DIR=/var/lib/scholartrace
Environment=SCHOLARTRACE_DB_PATH=/var/lib/scholartrace/scholartrace.db
ExecStart=/opt/miniconda3/envs/ScholarTrace/bin/scholartrace-mcp
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
LockPersonality=yes
RestrictSUIDSGID=yes
StateDirectory=scholartrace
ReadWritePaths=/var/lib/scholartrace
UMask=0077

[Install]
WantedBy=multi-user.target
